# ceu.timer
[Unit]
Description=Compress-Encrypt-Upload timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target 
